# systemd/envolees-cache-daily.timer
# Envolées - Timer de vérification cache QUOTIDIEN

[Unit]
Description=Envolées - Vérification cache quotidienne (timer)

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

# Pointer vers le service existant
Unit=envolees-cache.service

[Install]
WantedBy=timers.target
