# systemd/envolees-cache.service
# Envolées - Rafraîchissement du cache
#
# Maintient le cache à jour toutes les 6h
# Envoie une alerte si le cache est en erreur

[Unit]
Description=Envolées - Rafraîchissement cache
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/envolees

ExecStart=/usr/bin/bash -lc '\
    cd %h/dev/envolees && \
    source .venv/bin/activate && \
    set -a && source .env && source .env.secret 2>/dev/null || true && set +a && \
    python main.py cache-warm && \
    python main.py cache-verify || python main.py alert "Cache verify failed" --level warning \
'

# Timeout 10 min
TimeoutStartSec=600

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
