# systemd/envolees-validation-weekly.timer
# Envolées - Timer de validation HEBDOMADAIRE

[Unit]
Description=Envolées - Validation hebdomadaire (timer)

[Timer]
# Chaque dimanche à 20:00
OnCalendar=Sun *-*-* 20:00:00
Persistent=true
RandomizedDelaySec=300

# Pointer vers le service existant
Unit=envolees-validation.service

[Install]
WantedBy=timers.target
