# systemd/envolees-monthly.timer
# Envolées - Timer de réévaluation MENSUELLE
# Pipeline complet avec refresh forcé 1x/mois
#
# Le 1er de chaque mois à 19:00 :
# - Force le téléchargement de toutes les données fraîches
# - Recalcule tous les backtests IS/OOS
# - Génère une nouvelle shortlist
#
# C'est LA réévaluation importante pour ajuster ta liste d'instruments.

[Unit]
Description=Envolées - Réévaluation mensuelle complète (timer)

[Timer]
# Exécuter le 1er de chaque mois à 19:00
OnCalendar=*-*-01 19:00:00

Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
