# systemd/envolees-heartbeat.service
# Envolées - Heartbeat quotidien
#
# Envoie un signal de vie "tout va bien" 1x/jour
# Très simple, très rapide, très stable

[Unit]
Description=Envolées - Heartbeat quotidien
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/envolees

ExecStart=/usr/bin/bash -lc '\
    cd %h/dev/envolees && \
    source .venv/bin/activate && \
    set -a && source .env && source .env.secret 2>/dev/null || true && set +a && \
    python main.py heartbeat \
'

# Timeout court (30s max)
TimeoutStartSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
