# systemd/envolees-research.timer
# Envolées Research Timer
# Copier vers ~/.config/systemd/user/envolees-research.timer
#
# Déclenche la recherche 2x/jour : 07:30 et 19:30 (heure locale)

[Unit]
Description=Envolées research timer

[Timer]
# Déclenchement 2x/jour
OnCalendar=*-*-* 07:30:00
OnCalendar=*-*-* 19:30:00

# Rattraper les exécutions manquées (machine éteinte)
Persistent=true

# Délai aléatoire pour éviter les pics (optionnel)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
